# packages/nonla-default-settings/payload/etc/skel/.config/autostart/fcitx5.desktop
[Desktop Entry]
Type=Application
Name=Fcitx 5
Comment=Start Fcitx 5 input method
Exec=fcitx5 -d
Icon=fcitx
OnlyShowIn=KDE;
X-KDE-autostart-after=panel
X-KDE-StartupNotify=false
